# Ansible/templates/etc/profile.d/keepalivedenv.sh.j2
export KEEPALIVED_PRIORITY={{ KEEPALIVED_PRIORITY }}
